# deploy/systemd/ullm-openai.service.d/10-served-model.conf
[Service]
# Replace the legacy model environment instead of layering manifest mode over it.
EnvironmentFile=
EnvironmentFile=/etc/ullm/openai-gateway-manifest.env

# Preserve the base executable check and add fail-closed manifest validation.
ExecStartPre=
ExecStartPre=/usr/bin/test -x /home/homelab1/coding-local/ultimateLLM/uLLM-project/services/openai-gateway/.venv/bin/ullm-openai-gateway
ExecStartPre=/usr/bin/python3 /home/homelab1/coding-local/ultimateLLM/uLLM-project/tools/validate-served-model.py --manifest ${ULLM_SERVED_MODEL_MANIFEST}
